# docs/templates/local.Jenkins.Agent.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- launchd service name -->
    <key>Label</key>
    <string>local.Jenkins.Agent.launchd</string>

    <!--
      /bin/bash 経由で起動する。
      macOS のフルディスクアクセス/TCC対策として、run-agent.sh を直接実行しない。
    -->
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/your-user/jenkins-agent/run-agent.sh</string>
    </array>

    <!--
      launchd 自身の作業ディレクトリ。
      HDDやsymlink先は指定しない。
    -->
    <key>WorkingDirectory</key>
    <string>/Users/your-user/jenkins-agent</string>

    <!-- launchd は通常のログインシェル環境を読まないので明示 -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/your-user</string>

      <key>USER</key>
      <string>your-user</string>

      <key>LOGNAME</key>
      <string>your-user</string>

      <key>SHELL</key>
      <string>/bin/zsh</string>

      <key>LANG</key>
      <string>ja_JP.UTF-8</string>

      <key>LC_ALL</key>
      <string>ja_JP.UTF-8</string>

      <key>JAVA_HOME</key>
      <string>/opt/homebrew/opt/openjdk@21/libexec/openjdk.jdk/Contents/Home</string>

      <key>PATH</key>
      <string>/opt/homebrew/opt/openjdk@21/libexec/openjdk.jdk/Contents/Home/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- macOSログイン時に起動 -->
    <key>RunAtLoad</key>
    <true/>

    <!--
      agent が終了した場合に再起動。
      NetworkState は外す。ネットワーク状態変化で変な再起動を避ける。
    -->
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <!-- 再起動が過密にならないようにする -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- launchd 起動まわりのログ -->
    <key>StandardOutPath</key>
    <string>/Users/your-user/jenkins-agent/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/your-user/jenkins-agent/launchd.err.log</string>
  </dict>
</plist>
